# deploy/ip-tools-mcp.service
[Unit]
Description=ip-tools MCP Server
After=network.target

[Service]
Type=simple
User=ip-tools
Group=ip-tools
WorkingDirectory=/opt/ip-tools
EnvironmentFile=/opt/ip-tools/.env
ExecStart=/opt/ip-tools/.venv/bin/fastmcp run src/ip_tools/mcp/server.py:mcp --transport streamable-http --stateless --host 127.0.0.1 --port 8002
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ip-tools-mcp

[Install]
WantedBy=multi-user.target
